# LOCAL_BACKUP_v1.0.0_2025-10-12_104702/config/systemd/waxvalue-backend.service
[Unit]
Description=WaxValue Backend (FastAPI)
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=waxvalue
Group=waxvalue
WorkingDirectory=/home/waxvalue/waxvalue/backend
Environment=PYTHONPATH=/home/waxvalue/waxvalue/backend
EnvironmentFile=/home/waxvalue/waxvalue/.env.production
ExecStart=/home/waxvalue/waxvalue/backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --workers 4
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=waxvalue-backend

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/waxvalue/waxvalue/backend
ReadWritePaths=/var/log/waxvalue
ReadWritePaths=/home/waxvalue/backups

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
